# deployment/supervisor/berta.conf
[program:berta]
command=/home/pi/berta/venv/bin/gunicorn -b localhost:8000 -w 2 berta:app
directory=/home/pi/berta
user=pi
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
